# templates/migration-tools/supervisord.conf
[supervisord]
nodaemon = true
loglevel = info
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
logfile  = %(ENV_LOG_DIR)s/supervisord.log
pidfile  = %(ENV_SUP_RUN_DIR)s/supervisord.pid

[unix_http_server]
file = %(ENV_SUP_RUN_DIR)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix://%(ENV_SUP_RUN_DIR)s/supervisor.sock

[program:migration-worker]
command = migration-worker
stdout_logfile = %(ENV_LOG_DIR)s/%(program_name)s-%(process_num)s.log
stderr_logfile = %(ENV_LOG_DIR)s/%(program_name)s-%(process_num)s-stderr.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=7
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=7
autorestart = true
autostart = true

[program:migration-api-server]
command = migration-api-server
stdout_logfile = %(ENV_LOG_DIR)s/%(program_name)s-%(process_num)s.log
stderr_logfile = %(ENV_LOG_DIR)s/%(program_name)s-%(process_num)s-stderr.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=7
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=7
autorestart = true
autostart = true

[program:asynqmon]
command = asynqmon -redis-url 'redis://localhost:6379' -port 9102
stdout_logfile = %(ENV_LOG_DIR)s/%(program_name)s-%(process_num)s.log
stderr_logfile = %(ENV_LOG_DIR)s/%(program_name)s-%(process_num)s-stderr.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=7
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=7
autorestart = true
autostart = true
